# web/files/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=51200
logfile=/var/log/supervisor/supervisord-error.log
logfile_maxbytes=51200

[inet_http_server]
port=127.0.0.1:9001
username=supervisor
password=supervisor

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock ; use a unix:// URL for a unix socket 

[program:rsyslog]
command=/usr/sbin/rsyslogd -n
priority=1

[program:cron]
command=/entrypoint_cron.sh
autorestart=true
redirect_stderr=true
logfile=/var/log/cron.log
logfile_maxbytes=51200

[program:postfix]
process_name=master
directory=/etc/postfix
command=/usr/lib/postfix/sbin/master -c /etc/postfix -d
user=root
startsecs=0
autorestart=false

[program:redis-server]
command=/usr/bin/redis-server /etc/redis/redis.conf
user=redis

[program:apache2]
command=/bin/bash -c "source /etc/apache2/envvars && exec /usr/sbin/apache2 -D FOREGROUND"

[program:misp-modules]
command=sudo -u www-data /bin/bash -c "/var/www/MISP/venv/bin/misp-modules -l 127.0.0.1 -s"
startsecs = 0
autorestart = false

[program:clamd]
command=/bin/bash -c "/usr/sbin/clamd -c /etc/clamav/clamd.conf"
autorestart=true
